# battery.service
[Unit]
Description=Battery Warning

[Service]
ExecStart=/usr/local/bin/battery

[Install]
WantedBy=multi-user.target
